# deploy/mjuton-update.service
[Unit]
Description=Update MJUTON from the GitHub deployment release
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/mjuton-update
TimeoutStartSec=10min
